# alertmanager/supervisord.conf
[supervisord]
nodaemon=true

[program:alertmanager]
command=alertmanager --config.file=/etc/alertmanager/alertmanager.yml
